# memcached.target
[Unit]
Description=Memcached NoSQL key+value stores

[Install]
WantedBy=multi-user.target
